# rc.local
#!/bin/bash
sh /data/etc/dbus-smartmeter/reinstalllocal.sh

